# infra/systemd/analytic-worker-telegram.service
[Unit]
Description=Analytic System RQ Worker (telegram queue)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=fervuld
Group=fervuld
WorkingDirectory=/home/fervuld/prod/analytic-system
EnvironmentFile=/home/fervuld/prod/analytic-system/.env
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/fervuld/prod/analytic-system/backend
Environment=TELEGRAM_RQ_QUEUE_NAME=telegram
ExecStart=/home/fervuld/prod/analytic-system/.venv/bin/rq worker -u ${REDIS_URL} ${TELEGRAM_RQ_QUEUE_NAME}
Restart=always
RestartSec=5
TasksMax=512
MemoryHigh=1200M
MemoryMax=1600M
StandardOutput=journal
StandardError=inherit
TimeoutStopSec=45
KillSignal=SIGTERM
LimitNOFILE=65535
NoNewPrivileges=true
PrivateTmp=true
ReadWritePaths=/home/fervuld/prod/analytic-system /tmp

[Install]
WantedBy=multi-user.target
